# 60a274bd-067c-42f4-a185-f42245d95e99.TMP/modules/cache/ansible/roles/cache/templates/redis.init.j2
#!/bin/sh
#
# Redis init script for Debian/Ubuntu and CentOS
# Generated by Ansible
#
# chkconfig: - 85 15
# description: Redis is a persistent key-value database
#
### BEGIN INIT INFO
# Provides:          redis{{ redis_port | default('6379') }}
# Required-Start:    $network $local_fs $remote_fs
# Required-Stop:     $network $local_fs $remote_fs
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: Redis data structure server
# Description:       Redis is an open source, advanced key-value store.
### END INIT INFO

# Source function library
if [ -f /etc/rc.d/init.d/functions ]; then
    . /etc/rc.d/init.d/functions
fi

if [ -f /etc/init.d/functions ]; then
    . /etc/init.d/functions
fi

if [ -f /lib/lsb/init-functions ]; then
    . /lib/lsb/init-functions
fi

REDIS_PORT={{ redis_port | default('6379') }}
EXEC=/usr/bin/redis-server
CLIEXEC=/usr/bin/redis-cli
PIDFILE=/var/run/redis/${REDIS_PORT}.pid
CONF="/etc/redis/${REDIS_PORT}.conf"
REDISUSER=redis

case "$1" in
    start)
        if [ -f $PIDFILE ]
        then
                echo "$PIDFILE exists, process is already running or crashed"
        else
                echo "Starting Redis server..."
                if [ -n "$REDISUSER" ]
                then
                    su $REDISUSER -c "$EXEC $CONF"
                else
                    $EXEC $CONF
                fi
        fi
        ;;
    stop)
        if [ ! -f $PIDFILE ]
        then
                echo "$PIDFILE does not exist, process is not running"
        else
                PID=$(cat $PIDFILE)
                echo "Stopping Redis server..."
                $CLIEXEC -p $REDIS_PORT shutdown
                while [ -x /proc/${PID} ]
                do
                    echo "Waiting for Redis to shutdown ..."
                    sleep 1
                done
                echo "Redis stopped"
        fi
        ;;
    status)
        if [ ! -f $PIDFILE ]
        then
                echo "Redis is not running"
                exit 1
        else
                PID=$(cat $PIDFILE)
                if [ -x /proc/${PID} ]
                then
                    echo "Redis is running"
                    exit 0
                else
                    echo "Redis is not running but $PIDFILE exists"
                    exit 1
                fi
        fi
        ;;
    restart)
        $0 stop
        $0 start
        ;;
    *)
        echo "Please use start, stop, restart or status as first argument"
        ;;
esac