# 60a274bd-067c-42f4-a185-f42245d95e99.TMP/modules/cache/ansible/roles/cache/templates/sentinel.upstart.conf.j2
# Redis Sentinel Upstart configuration
# Generated by Ansible

description "Redis Sentinel"
author "Ansible"

start on runlevel [2345]
stop on runlevel [016]

respawn
respawn limit 10 5

setuid redis
setgid redis

limit nofile {{ redis_ulimit | default('65536') }} {{ redis_ulimit | default('65536') }}

expect fork

pre-start script
    mkdir -p /var/run/redis
    chown redis:redis /var/run/redis
end script

exec /usr/bin/redis-sentinel /etc/redis-sentinel/{{ redis_sentinel_port | default('26379') }}.conf