# 17295ca6-6b5f-42a4-9bd7-7fc04ae0429b.TMP/modules/fastapi-tutorial/ansible/roles/fastapi_tutorial/templates/fastapi-tutorial.service.j2
[Unit]
Description=FastAPI Tutorial Service
After=network.target postgresql.service

[Service]
Type=simple
User=root
WorkingDirectory={{ app_dir }}
Environment="PATH={{ app_dir }}/venv/bin"
ExecStart={{ app_dir }}/venv/bin/uvicorn app.main:app --host {{ app_host | default('0.0.0.0') }} --port {{ app_port | default(8000) }}
Restart=always

[Install]
WantedBy=multi-user.target