# b6c246e1-c249-4d54-8cbb-7faf6f495a8d.TMP/modules/cache/ansible/roles/cache/templates/sentinel.init.j2
#!/bin/sh
#
# Redis Sentinel init script
# Generated by Ansible
#
# chkconfig: 2345 95 20
# description: Redis Sentinel daemon
#
### BEGIN INIT INFO
# Provides: redis_sentinel_{{ item.port | default('26379') }}
# Required-Start: $network $local_fs $remote_fs
# Required-Stop: $network $local_fs $remote_fs
# Default-Start: 2 3 4 5
# Default-Stop: 0 1 6
# Short-Description: Redis Sentinel
# Description: Redis Sentinel monitors Redis instances
### END INIT INFO

# Source function library
[ -f /etc/rc.d/init.d/functions ] && . /etc/rc.d/init.d/functions
[ -f /etc/sysconfig/redis_sentinel_{{ item.port | default('26379') }} ] && . /etc/sysconfig/redis_sentinel_{{ item.port | default('26379') }}

# Source networking configuration
[ -f /etc/sysconfig/network ] && . /etc/sysconfig/network

# Check networking is up
[ "$NETWORKING" = "no" ] && exit 0

REDIS_USER={{ redis_user | default('redis') }}
REDIS_GROUP={{ redis_group | default('redis') }}
REDIS_HOME={{ redis_home | default('/var/lib/redis') }}
SENTINEL_PORT={{ item.port | default('26379') }}
REDIS_CONF_PATH={{ redis_conf_path | default('/etc/redis') }}
SENTINEL_CONF_FILE=$REDIS_CONF_PATH/sentinel_$SENTINEL_PORT.conf
REDIS_PID_PATH={{ redis_pid_path | default('/var/run/redis') }}
SENTINEL_PID_FILE=$REDIS_PID_PATH/sentinel_$SENTINEL_PORT.pid
SENTINEL_EXECUTABLE={{ redis_sentinel_executable | default('/usr/local/bin/redis-sentinel') }}

start() {
    [ -f $SENTINEL_CONF_FILE ] || exit 6
    [ -x $SENTINEL_EXECUTABLE ] || exit 5
    echo -n "Starting Redis Sentinel: "
    daemon --user $REDIS_USER "$SENTINEL_EXECUTABLE $SENTINEL_CONF_FILE"
    RETVAL=$?
    echo
    [ $RETVAL -eq 0 ] && touch /var/lock/subsys/redis_sentinel_$SENTINEL_PORT
    return $RETVAL
}

stop() {
    echo -n "Stopping Redis Sentinel: "
    killproc -p $SENTINEL_PID_FILE redis-sentinel
    RETVAL=$?
    echo
    [ $RETVAL -eq 0 ] && rm -f /var/lock/subsys/redis_sentinel_$SENTINEL_PORT
    return $RETVAL
}

restart() {
    stop
    start
}

case "$1" in
    start)
        start
        ;;
    stop)
        stop
        ;;
    status)
        status -p $SENTINEL_PID_FILE redis-sentinel
        RETVAL=$?
        ;;
    restart|force-reload)
        restart
        ;;
    condrestart|try-restart)
        [ -f /var/lock/subsys/redis_sentinel_$SENTINEL_PORT ] && restart || :
        ;;
    *)
        echo "Usage: $0 {start|stop|status|restart|condrestart|try-restart|force-reload}"
        exit 2
esac
exit $RETVAL