# 9b6c91b5-aa38-464f-aba6-fc11d333cc05.TMP/modules/cache/ansible/roles/cache/templates/redis.upstart.conf.j2
description "Redis Server {{ item.port }}"

start on runlevel [23]
stop on shutdown

pre-start script
    mkdir -p /var/run/redis
    chown {{ redis_user }}:{{ redis_group }} /var/run/redis
end script

exec sudo -u {{ redis_user }} /usr/bin/redis-server {{ redis_conf_dir }}/{{ item.port }}.conf

respawn