# 2a4a4362-e2dd-4a45-8f6c-fe321f187953.TMP/modules/fastapi-tutorial/ansible/roles/fastapi_tutorial/templates/fastapi-tutorial.service.j2
[Unit]
Description=FastAPI Tutorial Service
After=network.target postgresql.service

[Service]
Type=simple
User=root
WorkingDirectory={{ app_dir }}
Environment="PATH={{ app_dir }}/venv/bin"
ExecStart={{ app_dir }}/venv/bin/uvicorn app.main:app --host {{ app_host }} --port {{ app_port }}
Restart=always

[Install]
WantedBy=multi-user.target