# 85c3d3e2-f8a8-485d-a3b7-632c2d94be56.TMP/modules/cache/ansible/roles/cache/templates/redis.rcinit.j2
#!/bin/sh
#
# Redis init script for FreeBSD
# Generated by Ansible
#
# PROVIDE: redis_{{ redis_port }}
# REQUIRE: LOGIN
# KEYWORD: shutdown

. /etc/rc.subr

name="redis_{{ redis_port }}"
rcvar="${name}_enable"

load_rc_config $name

: ${redis_{{ redis_port }}_enable:="NO"}
: ${redis_{{ redis_port }}_user="{{ redis_user | default('redis') }}"}
: ${redis_{{ redis_port }}_group="{{ redis_group | default('redis') }}"}
: ${redis_{{ redis_port }}_config="/etc/redis/{{ redis_port }}.conf"}

command="{{ redis_install_dir | default('/usr/local') }}/bin/redis-server"
command_args="${redis_{{ redis_port }}_config}"
pidfile="{{ redis_pidfile }}"

start_cmd="redis_start"
stop_cmd="redis_stop"

redis_start() {
    echo "Starting ${name}."
    /usr/sbin/daemon -u ${redis_{{ redis_port }}_user} -p ${pidfile} ${command} ${command_args}
}

redis_stop() {
    echo "Stopping ${name}."
    {{ redis_install_dir | default('/usr/local') }}/bin/redis-cli -p {{ redis_port }} shutdown
}

run_rc_command "$1"