# 85c3d3e2-f8a8-485d-a3b7-632c2d94be56.TMP/modules/cache/ansible/roles/cache/templates/redis.upstart.conf.j2
# Redis Upstart script
# Generated by Ansible

description "Redis Server"
author "Ansible"

start on runlevel [2345]
stop on runlevel [016]

respawn
respawn limit 10 5

setuid {{ redis_user | default('redis') }}
setgid {{ redis_group | default('redis') }}

script
  exec {{ redis_install_dir | default('/usr/local') }}/bin/redis-server /etc/redis/{{ redis_port }}.conf
end script

pre-stop script
  exec {{ redis_install_dir | default('/usr/local') }}/bin/redis-cli -p {{ redis_port }} shutdown
end script