# 13f59d86-c088-4f4e-b08a-533ea5d92b6b.TMP/modules/cache/ansible/roles/cache/templates/redis.upstart.conf.j2
description "Redis Server {{ item.port }}"

start on runlevel [2345]
stop on runlevel [016]

respawn

setuid {{ redis_user | default('redis') }}
setgid {{ redis_group | default('redis') }}

exec /usr/local/bin/redis-server {{ redis_config_dir | default('/etc/redis') }}/{{ item.port }}.conf