# 9afd4b31-53e1-488c-8b86-56071f0c5386.TMP/ansible-project/roles/cache/templates/redis@.service.j2
[Unit]
Description=Redis persistent key-value database
After=network.target

[Service]
ExecStart=/usr/local/bin/redis-server /etc/redis/%i.conf
ExecStop=/usr/local/bin/redis-cli -p %i shutdown
User={{ user }}
Group={{ group }}
RuntimeDirectory=redis
RuntimeDirectoryMode=0755
LimitNOFILE={{ redis_maxclients | default('10000') }}

[Install]
WantedBy=multi-user.target