# 087f3e89-6dbe-4ebc-aa64-1d0cdddfac33.TMP/modules/fastapi-tutorial/ansible/roles/fastapi_tutorial/templates/fastapi-tutorial.service.j2
[Unit]
Description=FastAPI Tutorial Service
After=network.target postgresql.service

[Service]
Type=simple
User=root
WorkingDirectory={{ fastapi_app_dir | default('/opt/fastapi-tutorial') }}
Environment="PATH={{ fastapi_app_dir | default('/opt/fastapi-tutorial') }}/venv/bin"
ExecStart={{ fastapi_app_dir | default('/opt/fastapi-tutorial') }}/venv/bin/uvicorn app.main:app --host {{ fastapi_host | default('0.0.0.0') }} --port {{ fastapi_port | default('8000') }}
Restart=always

[Install]
WantedBy=multi-user.target